# lendingbot.service
[Unit]
Description=Polymarket Lending Bot (LIVE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/scannerpoly
ExecStart=/root/scannerpoly/.venv/bin/python bot_runner.py --live --window 360 --interval 15 --min-prob 0.50
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
